# deploy/systemd/weather-morning-report.service
[Unit]
Description=Send the daily weather morning report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=weather-report
Group=weather-report
WorkingDirectory=/opt/weather-morning-report
EnvironmentFile=/opt/weather-morning-report/.env
ExecStart=/opt/weather-morning-report/.venv/bin/weather-report send
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/weather-morning-report/var

[Install]
WantedBy=multi-user.target
