# deploy/systemd/weather-morning-report-worker.service
[Unit]
Description=Weather Morning Report v3 worker
After=network-online.target weather-morning-report-ui.service
Wants=network-online.target

[Service]
Type=simple
User=weather-report
Group=weather-report
WorkingDirectory=/opt/weather-morning-report
EnvironmentFile=/opt/weather-morning-report/.env
ExecStart=/opt/weather-morning-report/.venv/bin/weather-report serve-worker
Restart=on-failure
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/weather-morning-report/var

[Install]
WantedBy=multi-user.target
